# GB_HW_monitor.service
[Unit]
Description=GB_HW_Monitor

[Service]
ExecStart=/root/HW_monitor/GB_HW_monitor.sh

[Install]
WantedBy=multi-user.target
